# apps/rasp-deploy/src/garden.service
[Unit]
Description={{ service_name }}

[Service]
Type=simple
ExecStart={{ dest_service_bin }}
EnvironmentFile={{ dest_service_env }}
